# crontab
3 * * * * docker restart mackerel-plugin-aws-billing
